# Attack_Agent/result_handoff/APT32/APT32_composite_provenance_report.json
{
  "apt": "APT32",
  "variant_type": "composite",
  "dot_path": "D:\\notebook_UITNam3\\Nam3_ki2\\dacn\\TAGAPT\\TAGAPT_APTDetection\\Attack_Agent\\result_handoff\\APT32/APT32_composite_provenance.dot",
  "quality_valid": true,
  "stats": {
    "nodes": 42,
    "edges": 56,
    "clusters": 4
  },
  "reward": {
    "total": 0.7731771489300521,
    "coverage": 0.3333333333333333,
    "connectivity": 1.0,
    "validity": 1.0,
    "novelty": 0.5,
    "chain_depth": 0.9523809523809524,
    "transition_quality": 0.6960572035862341,
    "cluster_validity": 1.0,
    "macro_reward": 0.5186200212551819,
    "micro_reward": 0.9904761904761905
  },
  "planner_meta": {
    "variant_type": "composite",
    "composite_replacements": [
      {
        "replaced": "T1598.003",
        "replacement_sequence": [
          "T1598.002"
        ],
        "source_apt": "Star Blizzard",
        "similarity": 1.0
      },
      {
        "replaced": "T1589",
        "replacement_sequence": [
          "T1598.002"
        ],
        "source_apt": "Star Blizzard",
        "similarity": 1.0
      }
    ],
    "pruning": {
      "transitions": [
        {
          "source": "T1598.002",
          "target": "T1133",
          "score": 0.46491668963432314,
          "gnn_score": 0.7062121629714966,
          "scf_score": 0.17,
          "kept": false,
          "bridge_hint": "socket"
        },
        {
          "source": "T1133",
          "target": "T1078.003",
          "score": 0.844583169023196,
          "gnn_score": null,
          "scf_score": 0.844583169023196,
          "kept": true,
          "bridge_hint": "socket"
        },
        {
          "source": "T1078.003",
          "target": "T1053",
          "score": 0.778671752101183,
          "gnn_score": 0.7066076993942261,
          "scf_score": 0.8667500387430191,
          "kept": false,
          "bridge_hint": "socket"
        }
      ],
      "removed": [
        {
          "technique_id": "T1583.001",
          "reason": "weak transition and stronger/equal bypass",
          "incoming_score": 0.2333333333333333,
          "bypass_score": 0.46491668963432314
        }
      ]
    },
    "threshold": 0.8
  },
  "kill_chain": [
    {
      "technique_id": "T1598.002",
      "technique_name": "Spearphishing Attachment",
      "tactic": "reconnaissance",
      "context": {
        "technique_id": "T1598.002",
        "name": "Spearphishing Attachment",
        "tactics": [
          "reconnaissance"
        ],
        "platforms": [
          "pre"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "credential",
          "discovery",
          "file",
          "network"
        ],
        "requires_network": true,
        "requires_execution": false,
        "produces_access": true
      },
      "next_transition": {
        "source": "T1598.002",
        "target": "T1133",
        "score": 0.46491668963432314,
        "gnn_score": 0.7062121629714966,
        "scf_score": 0.17,
        "kept": false,
        "bridge_hint": "socket"
      }
    },
    {
      "technique_id": "T1133",
      "technique_name": "External Remote Services",
      "tactic": "persistence",
      "context": {
        "technique_id": "T1133",
        "name": "External Remote Services",
        "tactics": [
          "persistence",
          "initial-access"
        ],
        "platforms": [
          "containers",
          "linux",
          "macos",
          "windows"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "credential",
          "execution",
          "file",
          "network",
          "persistence",
          "process"
        ],
        "requires_network": true,
        "requires_execution": true,
        "produces_access": true
      },
      "next_transition": {
        "source": "T1133",
        "target": "T1078.003",
        "score": 0.844583169023196,
        "gnn_score": null,
        "scf_score": 0.844583169023196,
        "kept": true,
        "bridge_hint": "socket"
      }
    },
    {
      "technique_id": "T1078.003",
      "technique_name": "Local Accounts",
      "tactic": "stealth",
      "context": {
        "technique_id": "T1078.003",
        "name": "Local Accounts",
        "tactics": [
          "stealth",
          "persistence",
          "privilege-escalation",
          "initial-access"
        ],
        "platforms": [
          "containers",
          "esxi",
          "linux",
          "macos",
          "network",
          "windows"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "credential",
          "execution",
          "file",
          "network",
          "persistence",
          "privilege",
          "process"
        ],
        "requires_network": true,
        "requires_execution": true,
        "produces_access": true
      },
      "next_transition": {
        "source": "T1078.003",
        "target": "T1053",
        "score": 0.778671752101183,
        "gnn_score": 0.7066076993942261,
        "scf_score": 0.8667500387430191,
        "kept": false,
        "bridge_hint": "socket"
      }
    },
    {
      "technique_id": "T1053",
      "technique_name": "Scheduled Task/Job",
      "tactic": "execution",
      "context": {
        "technique_id": "T1053",
        "name": "Scheduled Task/Job",
        "tactics": [
          "execution",
          "persistence",
          "privilege-escalation"
        ],
        "platforms": [
          "containers",
          "esxi",
          "linux",
          "macos",
          "network",
          "windows"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "execution",
          "file",
          "network",
          "persistence",
          "privilege",
          "process"
        ],
        "requires_network": true,
        "requires_execution": true,
        "produces_access": true
      },
      "next_transition": null
    }
  ],
  "cluster_reports": [
    {
      "technique_id": "T1598.002",
      "tactic": "reconnaissance",
      "valid": true,
      "issues": [],
      "edge_types": [
        "EX",
        "FR",
        "RF",
        "WF"
      ],
      "node_types": [
        "file",
        "process"
      ]
    },
    {
      "technique_id": "T1133",
      "tactic": "persistence",
      "valid": true,
      "issues": [],
      "edge_types": [
        "EX",
        "FR",
        "RF",
        "WF"
      ],
      "node_types": [
        "file",
        "process"
      ]
    },
    {
      "technique_id": "T1078.003",
      "tactic": "stealth",
      "valid": true,
      "issues": [],
      "edge_types": [
        "EX",
        "FR",
        "RF",
        "WF"
      ],
      "node_types": [
        "file",
        "process"
      ]
    },
    {
      "technique_id": "T1053",
      "tactic": "execution",
      "valid": true,
      "issues": [],
      "edge_types": [
        "FR",
        "RF",
        "WF"
      ],
      "node_types": [
        "file",
        "process"
      ]
    }
  ]
}