# Attack_Agent/result_handoff/Lazarus/Lazarus_evolved_provenance_report.json
{
  "apt": "Lazarus",
  "variant_type": "evolved",
  "dot_path": "D:\\notebook_UITNam3\\Nam3_ki2\\dacn\\TAGAPT\\TAGAPT_APTDetection\\Attack_Agent\\result_handoff\\Lazarus/Lazarus_evolved_provenance.dot",
  "quality_valid": true,
  "stats": {
    "nodes": 43,
    "edges": 58,
    "clusters": 4
  },
  "reward": {
    "total": 0.8082510252346595,
    "coverage": 0.5,
    "connectivity": 1.0,
    "validity": 1.0,
    "novelty": 0.5,
    "chain_depth": 1.0,
    "transition_quality": 0.5825102523465951,
    "cluster_validity": 1.0,
    "macro_reward": 0.5538785883213083,
    "micro_reward": 1.0
  },
  "planner_meta": {
    "variant_type": "evolved",
    "composite_replacements": [],
    "pruning": {
      "transitions": [
        {
          "source": "T1589.002",
          "target": "T1053",
          "score": 0.46112826669216156,
          "gnn_score": 0.7095513939857483,
          "scf_score": 0.1575,
          "kept": false,
          "bridge_hint": "socket"
        },
        {
          "source": "T1053",
          "target": "T1037",
          "score": 0.6877137621343136,
          "gnn_score": 0.5402751564979553,
          "scf_score": 0.8679165023565293,
          "kept": false,
          "bridge_hint": "file"
        },
        {
          "source": "T1037",
          "target": "T1046",
          "score": 0.5986887282133102,
          "gnn_score": null,
          "scf_score": 0.5986887282133102,
          "kept": false,
          "bridge_hint": "socket"
        }
      ],
      "removed": [
        {
          "technique_id": "T1195",
          "reason": "weak transition and stronger/equal bypass",
          "incoming_score": 0.16799999999999998,
          "bypass_score": 0.46112826669216156
        }
      ]
    },
    "threshold": 0.8
  },
  "kill_chain": [
    {
      "technique_id": "T1589.002",
      "technique_name": "Email Addresses",
      "tactic": "reconnaissance",
      "context": {
        "technique_id": "T1589.002",
        "name": "Email Addresses",
        "tactics": [
          "reconnaissance"
        ],
        "platforms": [
          "pre"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "credential",
          "discovery",
          "file",
          "network"
        ],
        "requires_network": true,
        "requires_execution": false,
        "produces_access": true
      },
      "next_transition": {
        "source": "T1589.002",
        "target": "T1053",
        "score": 0.46112826669216156,
        "gnn_score": 0.7095513939857483,
        "scf_score": 0.1575,
        "kept": false,
        "bridge_hint": "socket"
      }
    },
    {
      "technique_id": "T1053",
      "technique_name": "Scheduled Task/Job",
      "tactic": "execution",
      "context": {
        "technique_id": "T1053",
        "name": "Scheduled Task/Job",
        "tactics": [
          "execution",
          "persistence",
          "privilege-escalation"
        ],
        "platforms": [
          "containers",
          "esxi",
          "linux",
          "macos",
          "network",
          "windows"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "execution",
          "file",
          "network",
          "persistence",
          "privilege",
          "process"
        ],
        "requires_network": true,
        "requires_execution": true,
        "produces_access": true
      },
      "next_transition": {
        "source": "T1053",
        "target": "T1037",
        "score": 0.6877137621343136,
        "gnn_score": 0.5402751564979553,
        "scf_score": 0.8679165023565293,
        "kept": false,
        "bridge_hint": "file"
      }
    },
    {
      "technique_id": "T1037",
      "technique_name": "Boot or Logon Initialization Scripts",
      "tactic": "persistence",
      "context": {
        "technique_id": "T1037",
        "name": "Boot or Logon Initialization Scripts",
        "tactics": [
          "persistence",
          "privilege-escalation"
        ],
        "platforms": [
          "esxi",
          "linux",
          "macos",
          "network",
          "windows"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "credential",
          "execution",
          "file",
          "network",
          "persistence",
          "privilege",
          "process"
        ],
        "requires_network": true,
        "requires_execution": true,
        "produces_access": true
      },
      "next_transition": {
        "source": "T1037",
        "target": "T1046",
        "score": 0.5986887282133102,
        "gnn_score": null,
        "scf_score": 0.5986887282133102,
        "kept": false,
        "bridge_hint": "socket"
      }
    },
    {
      "technique_id": "T1046",
      "technique_name": "Network Service Discovery",
      "tactic": "discovery",
      "context": {
        "technique_id": "T1046",
        "name": "Network Service Discovery",
        "tactics": [
          "discovery"
        ],
        "platforms": [
          "containers",
          "iaas",
          "linux",
          "macos",
          "network",
          "windows"
        ],
        "privilege_score": 0.0,
        "impact": false,
        "data_sources": [],
        "capabilities": [
          "discovery",
          "file",
          "network"
        ],
        "requires_network": true,
        "requires_execution": false,
        "produces_access": false
      },
      "next_transition": null
    }
  ],
  "cluster_reports": [
    {
      "technique_id": "T1589.002",
      "tactic": "reconnaissance",
      "valid": true,
      "issues": [],
      "edge_types": [
        "EX",
        "FR",
        "RF",
        "WF"
      ],
      "node_types": [
        "file",
        "process"
      ]
    },
    {
      "technique_id": "T1053",
      "tactic": "execution",
      "valid": true,
      "issues": [],
      "edge_types": [
        "FR",
        "RF",
        "WF"
      ],
      "node_types": [
        "file",
        "process"
      ]
    },
    {
      "technique_id": "T1037",
      "tactic": "persistence",
      "valid": true,
      "issues": [],
      "edge_types": [
        "EX",
        "FR",
        "RF",
        "WF"
      ],
      "node_types": [
        "file",
        "process"
      ]
    },
    {
      "technique_id": "T1046",
      "tactic": "discovery",
      "valid": true,
      "issues": [],
      "edge_types": [
        "FR",
        "RCV",
        "RF",
        "ST",
        "WF"
      ],
      "node_types": [
        "file",
        "process",
        "socket"
      ]
    }
  ]
}